# rpi/move.service
[Unit]
Description=Re-encode

[Service]
User=pi
ExecStart=/home/pi/move.sh

[Install]
WantedBy=multi-user.target
